# supervisor/mongod.conf
[program:mongod]
command=/usr/bin/mongod --config /etc/mongod.conf
stdout_logfile=/app/logs/mongod/mongod.log
stderr_logfile=/app/logs/mongod/mongod.log
autorestart=true
autostart=true
priority=110